# ops/eap-post-evening.timer
[Unit]
Description=vexp.me Evening Post Scheduler Timer

[Timer]
OnCalendar=*-*-* 14:30:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
