# ops/eap-post-morning.timer
[Unit]
Description=vexp.me Morning Post Scheduler Timer

[Timer]
OnCalendar=*-*-* 02:30:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
