# ops/eap-refresh-predictions.service
[Unit]
Description=vexp.me Daily Prediction Refresh

[Service]
Type=oneshot
User=mahimalam2400
WorkingDirectory=/home/mahimalam2400/Fifa_project
ExecStart=/home/mahimalam2400/Fifa_project/.venv/bin/python3 bot/refresh_predictions.py
EnvironmentFile=/home/mahimalam2400/Fifa_project/.env
StandardOutput=journal
StandardError=journal
